# templates/wsl2/ubuntu-22/etc/profile.d/wsl-keeagent.sh.j2
#!/bin/sh

export SSH_AUTH_SOCK="{{ wsl_keeagent_ssh_agent_socket }}"
